# gamescope-session/gamescope-session.service.d/override.conf
# Fix: gamescope-session boot loop / "Failed to connect to wayland socket: wayland-0"
#
# The stock unit (/usr/lib/systemd/user/gamescope-session.service) does
#   UnsetEnvironment=DISPLAY XAUTHORITY
# but NOT WAYLAND_DISPLAY. When a prior Plasma session leaks
# WAYLAND_DISPLAY=wayland-0 into the systemd --user manager environment,
# gamescope inherits it, auto-selects the NESTED wayland backend, fails to
# connect to a (nonexistent) parent compositor, exits 1, and plasmalogin
# relaunch-storms the session -> machine appears stuck at boot.
#
# Forcing WAYLAND_DISPLAY (and the X leftovers) out makes gamescope's backend
# auto-detection fall through to DRM/KMS, which is correct for the embedded
# game-mode session.
[Service]
UnsetEnvironment=DISPLAY XAUTHORITY WAYLAND_DISPLAY
